# lamp/supervisord-apache2.conf
[program:apache2]
command=/opt/pra/start-apache2.sh
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
